# php7.0/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
echo "[services.d] started nginx"
# Run nginx processes as wordpress:web
nginx
