# php5.6/system-root/etc/services.d/nginx/run
#!/usr/bin/execlineb -P

# Run nginx processes as wordpress:web
nginx
