# system-root/etc/services.d/crond/run
#!/usr/bin/with-contenv sh

# Runs cron daemon in foreground
crond -l 2 -f
