# php7.0/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash
echo "[services.d] started php-fpm7"
exec /usr/bin/php-fpm7 -c /etc/php7/php-fpm.ini --fpm-config /etc/php7/php-fpm.conf --pid /var/run/php-fpm.pid

