# php7.0/rootfs/etc/services.d/crond/run
#!/usr/bin/with-contenv sh
echo "[services.d] started cron daemon"
# Runs cron daemon in foreground
crond -l 2 -f
