# php5.6/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash

exec /usr/bin/php-fpm -c /etc/php/php-fpm.ini --fpm-config /etc/php/php-fpm.conf --pid /var/run/php-fpm.pid

