# debian-php7.0/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

# Run nginx always in foreground
nginx -g "daemon off;"
